# microservice/image/supervisor/conf.d/image.conf
[supervisord]
nodaemon=true

[program:image]

# --error-logfile /workspace/gunicorn-error.log 
command = python -m app.main
directory = /workspace
autostart=true
autorestart=true
stdout_logfile=/workspace/logs/image_out.log
stderr_logfile=/workspace/logs/image_err.log