# microservice/classify/supervisor/conf.d/classify.conf
[supervisord]
nodaemon=true

[program:classify]

# --error-logfile /workspace/gunicorn-error.log 
command = python -m app.main
directory = /workspace
autostart=true
autorestart=true
stdout_logfile=/workspace/logs/classify_out.log
stderr_logfile=/workspace/logs/classify_err.log