# microservice/qafull/supervisor/conf.d/qa_full.conf
[supervisord]
nodaemon=true

[program:qafull]

command = python -m app.main
directory = /workspace
autostart=true
autorestart=true
stdout_logfile=/workspace/logs/qa_full_out.log
stderr_logfile=/workspace/logs/qa_full_err.log